# pkg/scripts/postinstall
#! /bin/sh

launchctl load /Library/LaunchAgents/adderall.plist
